# lib/headroom-proxy.service.in
[Unit]
Description=Headroom compression proxy (OpenRouter backend) for Hermes Agent
Documentation=https://github.com/chopratejas/headroom
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# OPENROUTER_API_KEY (and the rest of the project's .env) is loaded from here.
# HEADROOM_MODE (token|cache) is also read from here natively by the proxy.
EnvironmentFile=__ENVFILE__
ExecStart=__HEADROOM_BIN__ proxy --backend openrouter --host 127.0.0.1 --port __PORT__ $HEADROOM_EXTRA_ARGS
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
